# systemd/crypto-notice.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/root/main

[Install]
WantedBy=multi-user.target